# stock-report-site/com.william.stockreport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.william.stockreport</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/mavicair2tw/.openclaw/workspace/stock-report-site && ./.venv/bin/python scripts/update_report.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/mavicair2tw/.openclaw/workspace/stock-report-site/logs/stockreport.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mavicair2tw/.openclaw/workspace/stock-report-site/logs/stockreport.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
